# backend/robot-four-final.service
[Unit]
Description=Post-processing daemon

[Service]
Type=exec
ExecStart=/opt/robot-four/robot-four-final.py

[Install]
WantedBy=multi-user.target
